# deploy/photobooth-kiosk.service
[Unit]
Description=Photobooth Kiosk Display
After=graphical-session.target
Wants=photobooth.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStartPre=/bin/sh -c 'xset s off -dpms 2>/dev/null || true'
ExecStart=/bin/sh -c 'exec $(which chromium 2>/dev/null || which chromium-browser) \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-component-update \
    --check-for-update-interval=31536000 \
    --autoplay-policy=no-user-gesture-required \
    http://localhost:8000/booth'
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
